# deploy/gateway-be.service
[Unit]
Description=Gateway Backend Service
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=vpsroot
WorkingDirectory=/home/vpsroot/projects/gateway-be
EnvironmentFile=/home/vpsroot/projects/gateway-be/.env
ExecStart=/home/vpsroot/projects/gateway-be/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
